# additions/etc/periodic/weekly/roothints
#!/bin/sh
curl --dns-servers 127.0.0.1,[::1] -L -s -o /etc/unbound/root.hints https://www.internic.net/domain/named.root
service unbound reload
